# sync.service
[Unit]
Description=Sync neovim config to GitHub

[Service]
Type=oneshot
ExecStart=%h/projects/open-source/my-neovim/sync.sh
