# sync.timer
[Unit]
Description=Sync neovim config to GitHub weekly

[Timer]
OnCalendar=Sun *-*-* 20:00:00
Persistent=true

[Install]
WantedBy=timers.target
